# superpipes/src/deb/control/postinst
#!/bin/sh
# postinst script for superpipes
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

case "$1" in

  configure)
    /usr/sbin/groupadd -r superpipes || true
    /usr/sbin/useradd -d /usr/local/superpipes -g superpipes -M -r -s /bin/sh superpipes || true

    /bin/mkdir /usr/local/superpipes/log || true
    /bin/chmod 755 /usr/local/superpipes/log || true
    /usr/bin/touch /usr/local/superpipes/log/info || true
    /usr/bin/touch /usr/local/superpipes/log/error || true

    /usr/bin/touch /usr/local/superpipes/superpipes.pid || true

    /bin/mkdir /usr/local/superpipes/datas || true
    /bin/chmod 755 /usr/local/superpipes/datas || true

    /bin/chown superpipes:superpipes /usr/local/superpipes -Rf || true

    /usr/sbin/update-rc.d superpipes defaults || true
  ;;

  abort-upgrade|abort-remove|abort-deconfigure)
    exit 0
  ;;

  *)
    echo "postinst called with unknown argument \`$1'" >&2
    exit 1
  ;;

esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
